# axolotl/tests/data/docker/summing/etc/service/summing/run
#!/bin/bash -e

cd /code
exec chpst -u nobody:nogroup ./server.py 2>&1
